# conf/etc/service/grafana/run
#!/bin/bash

exec /opt/grafana/bin/grafana-server -config /opt/grafana/conf/sample.ini -homepath /opt/grafana/ 2>&1 >> /var/log/grafana.log

